# deb/backup-db/DEBIAN/prerm
#!/bin/bash

set -e

service backup-db stop
systemctl disable backup-db.service

exit 0
